# resources/cosmos-power-stream-host.service
# /etc/systemd/system/cosmos-power-stream-host.service
[Unit]
Description=Cosmos Power Stream Host
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/cosmos-power-stream/packages/backend
ExecStart=/usr/bin/pnpm run host
Restart=on-failure
Environment=POSTGRES_USER=cosmos

[Install]
WantedBy=multi-user.target
